# zone-node/systemd/zone-node-hotspot.service
[Unit]
Description=Zone Node Setup Hotspot (onboarding)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zone-node-hotspot-onboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
